# recipes-swupdate/swupdate/swupdate/swupdate-usb.rules
ACTION=="add", KERNEL=="sd*", SUBSYSTEM=="block", ENV{ID_FS_USAGE}=="filesystem", RUN+="/bin/systemctl --no-block start swupdate-usb@%k.service"
