# recipes-swupdate/swupdate/swupdate/swupdate-usb@.service
[Unit]
Description=usb media swupdate service

[Service]
Type=oneshot
ExecStartPre=
ExecStart=/usr/bin/swupdate-client %I
ExecStopPost=
